# wall/systemd_bpq-wall.socket.example
[Unit]
Description=BPQ Wall Socket

[Socket]
ListenStream=6789
Accept=yes

[Install]
WantedBy=sockets.target
